# supervisord.conf
[program:celeryd]
command = python /cnhome/manage.py celery worker -l info
stdout_logfile=/cnhome/logs/celery_out.log
stderr_logfile=/cnhome/logs/celery_err.log
directory=/cnhome/

